# apps/agent/systemd/noderail-agent.service
[Unit]
Description=NodeRail node agent
Documentation=https://noderail.dev/docs/agent
# The agent talks to Docker; start after the daemon is up.
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=noderail
Group=noderail
# The noderail user must be a member of the `docker` group to reach the socket.
# The install script adds it: `usermod -aG docker noderail`.
SupplementaryGroups=docker

ExecStart=/usr/local/bin/noderail-agent run --config /etc/noderail/agent.toml
Restart=always
RestartSec=5
# Give in-flight commands a moment to report results before SIGKILL.
TimeoutStopSec=30
KillSignal=SIGTERM

# State directory: systemd creates /var/lib/noderail owned by the service user.
StateDirectory=noderail
StateDirectoryMode=0750
WorkingDirectory=/var/lib/noderail

# Environment (optional): raise log verbosity with NODERAIL_LOG=debug.
Environment=NODERAIL_LOG=info

# --- Hardening ---------------------------------------------------------------
# The agent needs the Docker socket and outbound network; otherwise it is locked
# down as much as is compatible with driving Docker.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# Allow writing only to its state and config dirs.
ReadWritePaths=/var/lib/noderail /etc/noderail
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=false
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
