# apps/agent/systemd/yourstack-agent.service
[Unit]
Description=YourStack node agent
Documentation=https://yourstack.dev/docs/agent
# The agent talks to a container engine (Docker or Podman). Start after whichever
# is present; neither is hard-required so the unit works on either kind of host.
After=network-online.target docker.service podman.socket
Wants=network-online.target docker.service podman.socket

[Service]
Type=simple
User=yourstack
Group=yourstack
# The yourstack user reaches the engine socket via its group membership, which
# the install script grants: `usermod -aG docker|podman yourstack`. systemd
# initializes the service's supplementary groups from the user database, so no
# explicit SupplementaryGroups= (which would fail if the group is absent).

ExecStart=/usr/local/bin/yourstack-agent run --config /etc/yourstack/agent.toml
Restart=always
RestartSec=5
# Give in-flight commands a moment to report results before SIGKILL.
TimeoutStopSec=30
KillSignal=SIGTERM

# State directory: systemd creates /var/lib/yourstack owned by the service user.
StateDirectory=yourstack
StateDirectoryMode=0750
WorkingDirectory=/var/lib/yourstack

# Environment (optional): raise log verbosity with YOURSTACK_LOG=debug.
Environment=YOURSTACK_LOG=info

# --- Hardening ---------------------------------------------------------------
# The agent needs the engine socket and outbound network; otherwise it is locked
# down as much as is compatible with driving containers.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# Allow writing only to its state and config dirs.
ReadWritePaths=/var/lib/yourstack /etc/yourstack
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=false
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
